# source/application/CLSSR570Coordinator/clsSR570Coordinator
# Linux rc script for CLS SR status building display system
#
# chkconfig: 345 90 10
# description: CLS SR status building display system. Must be started after networking and NFS
#

INSTALL_DIR={IOCAPPDIR}/beamline/acquaman/CLSSR570Coordinator
TOP=/home/epics/src/R3.14.12/
T_A=linux-x86_64
RUNSCRIPT=runCLSSR570CoordinatorStartup
DAEMON=CLSAddOnsCoordinator
PORTNO={PORTNO}
NETWORK="{NETWORK}"

source /iocApps/Common/startup.include

